# deploy/supervisord.conf
[supervisord]
logfile=/data/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
childlogdir=/data/log/

[inet_http_server]
port=0.0.0.0:9005
username=nomooc              ;登录管理后台的用户名
password=nomooc^*               ;登录管理后台的密码

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://0.0.0.0:9005
username=glsttome              ;登录管理后台的用户名
password=glst123               ;登录管理后台的密码

[program:nginx]
command=nginx -c /app/deploy/nginx/nginx.conf
directory=/app/
stdout_logfile=/data/log/nginx.log
stderr_logfile=/data/log/nginx.log
autostart=true
autorestart=true
startsecs=3
stopwaitsecs = 10
killasgroup=true

[program:gunicorn]
command=gunicorn oj.wsgi --user server --group spj --bind 127.0.0.1:8080 --workers 5 --threads 4 --max-requests-jitter 10000 --max-requests 100000 --keep-alive 32 -k gevent
directory=/app/
stdout_logfile=/data/log/gunicorn.log
stderr_logfile=/data/log/gunicorn.log
autostart=true
autorestart=true
startsecs=3
stopwaitsecs = 10
killasgroup=true

[program:celery]
command=celery -A oj worker -l warning --autoscale 2,4
directory=/app/
user=nobody
stdout_logfile=/data/log/celery.log
stderr_logfile=/data/log/celery.log
autostart=true
autorestart=true
startsecs=3
stopwaitsecs = 10
killasgroup=true

[program:beat_celery]
command=celery -A oj beat -l warning
directory=/app/
stdout_logfile=/data/log/celery_beat.log
stderr_logfile=/data/log/celery_beat.log
autostart=true
autorestart=true
startsecs=3
stopwaitsecs = 10
killasgroup=true

[program:monitor_celery]
command=celery flower -A oj --port=8088  --basic_auth=nomooc:nomooc^*
stdout_logfile=/data/log/celery_monitor.log
stderr_logfile=/data/log/celery_monitor.log
autostart=true
autorestart=true
startsecs=3
stopwaitsecs = 10
killasgroup=true
